# cmdeploy/src/cmdeploy/systemd/chatmail-oauth2.service
[Unit]
Description=Chatmail OAuth2 Authentication Service
After=network.target

[Service]
Type=simple
User=vmail
Group=vmail
WorkingDirectory=/usr/local/lib/chatmaild
ExecStart=/usr/local/lib/chatmaild/venv/bin/python /usr/local/lib/chatmaild/venv/bin/chatmail-oauth2 /usr/local/lib/chatmaild/chatmail.ini
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/vmail

[Install]
WantedBy=multi-user.target
